# ecommerce-api.service
[Unit]
Description=E-Commerce FastAPI Backend Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=ecommerce
WorkingDirectory=/opt/ecommerce-api
Environment="PATH=/opt/ecommerce-api/venv/bin"
EnvironmentFile=/opt/ecommerce-api/.env
ExecStart=/opt/ecommerce-api/venv/bin/gunicorn \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --access-logfile /var/log/ecommerce-api/access.log \
    --error-logfile /var/log/ecommerce-api/error.log \
    --log-level info \
    main:app

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ecommerce-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
